# esp-devices/files/esp_devices_script
#!/bin/sh /etc/rc.common

START=99
STOP=1
 
start() {        
    /etc/init.d/esp_devices start
}                 
 
stop() {          
    /etc/init.d/esp_devices stop
}